# solr/templates/default/solr-jetty.erb
#!/bin/bash
### BEGIN INIT INFO
# Provides:          <%= ["solr", @role].compact.join("_") %>
# Required-Start:    $remote_fs $network $syslog
# Required-Stop:     $remote_fs $network $syslog
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Solr
# Description:       Solr
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=solr
PIDFILE=/var/run/<%= ["solr", @role].compact.join("_") %>.pid

JAVA_HEAP_SIZE="<%= node.solr.java_heap_size %>"
JAVA_NEW_SIZE="<%= node.solr.java_new_size %>"
JAVA_SURVIVOR_RATIO="<%= node.solr.java_survivor_ratio %>"

# Garbage collections - use concurrent low pause collectors
JAVA_OPTS="-XX:+UseParallelGC -XX:+PrintGCDetails -XX:+PrintGCTimeStamps"

# Memory settings: Heap size
JAVA_OPTS="$JAVA_OPTS -Xmx$JAVA_HEAP_SIZE -Xms$JAVA_NEW_SIZE +PrintCommandLineFlags"

#  $JAVA_OPTS $NEWRELIC_OPTS
START_PREFIX="start-stop-daemon -c $USER -m --pidfile $PIDFILE --start --chdir /opt/apache-solr/example --exec /usr/bin/java -- $JAVA_OPTS"

SOLR_HOME="<%= @solr_home %>"
if [[ "$PREFILL" ]]; then
  SOLR_PORT="<%= @solr_port.to_i + 1000 %>"
else
  SOLR_PORT="<%= @solr_port %>"
fi
SOLR_OPTS="-Dsolr.solr.home=$SOLR_HOME -Djetty.port=$SOLR_PORT"
<% if @solr_host %>
SOLR_OPTS="$SOLR_OPTS -Djetty.host=<%= @solr_host %>"
<% end %>
SOLR_MASTER_URL="<%= node.solr.master_url %>"

set -e

# Parse command line parameters.
case $1 in
    start)
        echo "starting SOLR"
        $($START_PREFIX $SOLR_OPTS <%= @master_slave_options %> -jar start.jar 2>&1 | logger -t solr &)
        ;;
    stop)
        echo "stopping SOLR"
        start-stop-daemon --pidfile $PIDFILE --stop
        ;;
    restart)
  $0 stop
  sleep 10s
  $0 start
        ;;
    *)
        # Print help
        echo "Usage: $0 {start|start_master|start_slave|    stop|restart}" 1>&2
        exit 1
        ;;
esac

exit 0
