# jenkins/templates/default/jenkins_start_stop.erb
#!/bin/bash
### BEGIN INIT INFO
# Provides:          jenkins
# Required-Start:    $remote_fs $network $syslog
# Required-Stop:     $remote_fs $network $syslog
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Jenkins
# Description:       Jenkins
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=<%= node.jenkins.user %>
PIDFILE=/var/run/jenkins.pid

JENKINS_OPTS="--httpPort=<%= node.jenkins.port %> --httpListenAddress=<%= node.jenkins.address %>"
<% if node.jenkins[:prefix] %>
JENKINS_OPTS="$JENKINS_OPTS --prefix=<%= node.jenkins[:prefix] %>"
<% end %>

START_PREFIX="start-stop-daemon -c $USER -m --pidfile $PIDFILE --start --chdir /home/$USER --exec /usr/bin/java -- "

set -e

# Parse command line parameters.
case $1 in
    start)
        echo "starting jenkins"
        $($START_PREFIX -jar /opt/jenkins/lib/jenkins.war $JENKINS_OPTS 2>&1 | logger -t jenkins &)
        ;;
    stop)
        echo "stopping jenkins"
        start-stop-daemon --pidfile $PIDFILE --stop
        ;;
    restart)
  $0 stop
  sleep 10s
  $0 start
        ;;
    *)
        # Print help
        echo "Usage: $0 {start|start_master|start_slave|    stop|restart}" 1>&2
        exit 1
        ;;
esac

exit 0
